# services/ecowatch_damy.service
[Unit]
Description = ecowatch service

[Service]
ExecStart   = sh /home/procon/ecowatch_project/scripts/start_with_damy.sh
Type        = simple

[Install]
WantedBy    = default.target